# system/blur-power/90-blur-power.rules
# Sync Hyprland blur to the power source on every AC plug/unplug.
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", RUN+="/usr/local/bin/blur-power-sync"
