# system/anti-heat/anti-heat-sync.service
[Unit]
Description=Sync the anti-heat undervolt to the caelestia toggle state

[Service]
Type=oneshot
ExecStart=/usr/local/bin/anti-heat-sync

[Install]
WantedBy=multi-user.target
